# code/temp_monitor.cron
* * * * * /usr/bin/python Embedded_Linux/ELSpring2018/code/cpu_temp_monitor.py
