# WebAssignment/tempLog.cron
0,15,30,45 * * * *  /usr/bin/bin/python /home/pi/Embedded_Linux/WebAssignment/rT.py
